# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:desktop]
command=/usr/local/bin/start-vnc.sh
environment=HOME="/root",USER="root"

autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

priority=10